# qyk3sworker.service
[Unit]
Description=K3s-worker
After=network.target

[Service]
Type=simple
ExecStart=/root/startasworker2.sh

[Install]
WantedBy=multi-user.target


