# qyk3server.service
[Unit]
Description=K3s
After=network.target

[Service]
Type=simple
ExecStart=/root/k3s server --docker

[Install]
WantedBy=multi-user.target


